# supervisor/cinder-volume.ini
[program:cinder-volume]
user=cinder
command=/usr/bin/cinder-volume --config-file /etc/cinder/cinder.conf --logfile /var/log/cinder/volume.log
priority=200
autostart=false
autorestart=true
stderr_logfile=/var/log/supervisor/cinder-volume.err.log
stdout_logfile=/var/log/supervisor/cinder-volume.out.log
