# supervisor/httpd.ini
[program:httpd]
user=root
command=/usr/sbin/httpd -DFOREGROUND
priority=1
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/httpd.err.log
stdout_logfile=/var/log/supervisor/httpd.out.log
