# supervisor/mariadb.ini
[program:mariadb]
user=mysql
command=/usr/local/bin/mariadb.sh
priority=1
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mariadb.err.log
stdout_logfile=/var/log/supervisor/mariadb.out.log
