# supervisor/prepare.ini
[program:prepare]
user=root
command=/usr/local/bin/prepare.sh
priority=200
autostart=true
autorestart=false
stderr_logfile=/var/log/supervisor/prepare.err.log
stdout_logfile=/var/log/supervisor/prepare.out.log
