# supervisor/rabbitmq.ini
[program:rabbitmq]
user=root
command=/usr/local/bin/rabbitmq.sh
priority=10
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/rabbitmq.err.log
stdout_logfile=/var/log/supervisor/rabbitmq.out.log
